# dockerfile/dockerfile/fdommp-plugins/app/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock


[program:websocket]
command=/usr/local/bin/daphne -p 8090 -b 0.0.0.0 fdommp.asgi:application
directory=/mnt/fdommp
stdout_logfile=/mnt/fdommp/logs/fdommp-websocket.log
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
numprocs=1



[program:celery-worker-default]
command=/usr/local/bin/celery -A fdommp  worker --loglevel=info
directory=/mnt/fdommp
stdout_logfile=/mnt/fdommp/logs/fdommp-celery-worker-default.log
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
numprocs=1

